# systemd/pskhouse.target
[Unit]
Description=PSKHouse
After=multi-user.target network.target

[Install]
WantedBy=multi-user.target
